# supervisord-postgresql.conf
[program:mariadb]
command=/supervisord-postgresql.sh
numprocs=1
autostart=true
autorestart=true